# package/upstart.conf
description "The Service"

start on runlevel [2345]
stop on shutdown
respawn

script
  `get-config` && node /opt/the-app/app.js
end script
