# package/upstart.conf
description "The Service"

start on runlevel [2345]
stop on shutdown
respawn

exec node /opt/the-app/app.js
